# addons/actual-budget-complete/rootfs/etc/services.d/actual-budget-addon/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Actual Budget Complete Add-on - Finish Script
# ==============================================================================

bashio::log.info "🛑 Actual Budget add-on shutting down..."

# Clean up any resources if needed
bashio::log.info "✅ Shutdown complete"