# addons/actual-budget-complete/rootfs/etc/services.d/actual-budget-addon/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Actual Budget Complete Add-on - Native HA Integration
# Creates dashboard tiles and entities for Home Assistant
# ==============================================================================

# Wait for Home Assistant API to be available
bashio::log.info "ðŸ  Starting Actual Budget native Home Assistant integration..."
bashio::log.info "â³ Waiting for Home Assistant API..."

# Wait for supervisor token
while [[ -z "${SUPERVISOR_TOKEN}" ]]; do
    sleep 1
done

bashio::log.info "âœ… Home Assistant API ready"

# Parse configuration
MANAGE_API_SERVER=$(bashio::config 'manage_api_server')
MANAGE_ACTUAL_SERVER=$(bashio::config 'manage_actual_server')
API_SERVER_PORT=$(bashio::config 'api_server_port')
ACTUAL_SERVER_PORT=$(bashio::config 'actual_server_port')
LOG_LEVEL=$(bashio::config 'log_level')

# Set log level
bashio::log.level "${LOG_LEVEL}"

bashio::log.info "ðŸ“Š Creating native Home Assistant dashboard entities..."

# Function to create entity
create_entity() {
    local entity_id="$1"
    local name="$2"
    local icon="$3"
    local initial_state="$4"
    local attributes="$5"
    
    local payload="{
        \"state\": \"${initial_state}\",
        \"attributes\": {
            \"friendly_name\": \"${name}\",
            \"icon\": \"${icon}\",
            ${attributes}
        }
    }"
    
    curl -s -X POST \
        -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
        -H "Content-Type: application/json" \
        -d "${payload}" \
        "http://supervisor/core/api/states/${entity_id}" > /dev/null 2>&1
    
    if [[ $? -eq 0 ]]; then
        bashio::log.debug "âœ… Created entity: ${entity_id}"
    else
        bashio::log.warning "âš ï¸ Failed to create entity: ${entity_id}"
    fi
}

# Create dashboard entities
bashio::log.info "ðŸŽ¯ Creating status tiles..."

create_entity "sensor.actual_budget_system_status" \
    "System Status" \
    "mdi:calculator-variant" \
    "running" \
    '"description": "Overall system health", "tooltip": "Green=Running, Yellow=Issues, Red=Failed", "version": "2.0.0"'

create_entity "sensor.actual_budget_api_status" \
    "API Server" \
    "mdi:api" \
    "external" \
    '"description": "Node.js API server status", "tooltip": "Connect to your existing API server or let add-on manage it", "managed": false'

create_entity "sensor.actual_budget_server_status" \
    "Actual Budget Server" \
    "mdi:server" \
    "external" \
    '"description": "Official Actual Budget server", "tooltip": "Connect to your existing server or let add-on manage it", "managed": false'

create_entity "sensor.actual_budget_last_sync" \
    "Last Sync" \
    "mdi:sync" \
    "never" \
    '"description": "Last synchronization time", "tooltip": "When transactions were last synced to Xero"'

bashio::log.info "ðŸ“ˆ Creating statistics tiles..."

create_entity "sensor.actual_budget_transactions_processed" \
    "Transactions Processed" \
    "mdi:counter" \
    "0" \
    '"unit_of_measurement": "transactions", "description": "Total processed", "tooltip": "ðŸ”¸ PLACEHOLDER: Shows test data until first sync completes"'

create_entity "sensor.actual_budget_successful_imports" \
    "Successful Imports" \
    "mdi:check-circle" \
    "0" \
    '"unit_of_measurement": "transactions", "description": "Successfully imported to Xero", "tooltip": "ðŸ”¸ PLACEHOLDER: Shows test data until first sync completes"'

create_entity "sensor.actual_budget_failed_transactions" \
    "Failed Transactions" \
    "mdi:alert-circle" \
    "0" \
    '"unit_of_measurement": "transactions", "description": "Failed imports", "tooltip": "ðŸ”¸ PLACEHOLDER: Shows test data until first sync completes"'

create_entity "sensor.actual_budget_pending_mappings" \
    "Pending Mappings" \
    "mdi:map-marker-question" \
    "0" \
    '"unit_of_measurement": "mappings", "description": "Awaiting category mapping", "tooltip": "ðŸ”¸ PLACEHOLDER: Shows test data until first sync completes"'

bashio::log.info "ðŸŽ‰ Native HA dashboard integration complete!"
bashio::log.info "ðŸ“± Look for 'Actual Budget' in your Home Assistant sidebar"
bashio::log.info "ðŸ“Š Dashboard tiles created with placeholder data"
bashio::log.info "ðŸ’¡ Configure external service URLs to connect to your existing servers"

# Main service loop - keep the add-on running and update entities periodically
while true; do
    # Update system status with current timestamp
    curl -s -X POST \
        -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
        -H "Content-Type: application/json" \
        -d "{
            \"state\": \"running\",
            \"attributes\": {
                \"friendly_name\": \"System Status\",
                \"icon\": \"mdi:calculator-variant\",
                \"description\": \"Native HA integration active\",
                \"last_update\": \"$(date -Iseconds)\",
                \"version\": \"2.0.0\",
                \"integration_type\": \"native_dashboard\"
            }
        }" \
        "http://supervisor/core/api/states/sensor.actual_budget_system_status" > /dev/null 2>&1
    
    # Sleep for 5 minutes before next update
    sleep 300
done