# addons/actual-budget-complete/rootfs/etc/services.d/actual-budget-addon/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Actual Budget Complete Add-on - Native HA Integration
# Creates dashboard tiles and entities for Home Assistant
# ==============================================================================

# Enable error handling but don't exit on errors
set +e

# Wait for Home Assistant API to be available
bashio::log.info "🏠 Starting Actual Budget native Home Assistant integration..."
bashio::log.info "⏳ Waiting for Home Assistant API..."

# Wait for supervisor token with timeout
timeout=60
counter=0
while [[ -z "${SUPERVISOR_TOKEN}" && $counter -lt $timeout ]]; do
    sleep 1
    ((counter++))
done

if [[ -z "${SUPERVISOR_TOKEN}" ]]; then
    bashio::log.error "❌ Failed to get supervisor token after ${timeout} seconds"
    exit 1
fi

bashio::log.info "✅ Home Assistant API ready"

# Parse configuration with defaults
EXTERNAL_API_URL=$(bashio::config 'external_api_url' '')
EXTERNAL_ACTUAL_URL=$(bashio::config 'external_actual_url' '')
API_SERVER_PORT=$(bashio::config 'api_server_port' '3001')
ACTUAL_SERVER_PORT=$(bashio::config 'actual_server_port' '5006')
LOG_LEVEL=$(bashio::config 'log_level' 'info')

# Set log level
bashio::log.level "${LOG_LEVEL}"

bashio::log.info "📊 Creating native Home Assistant dashboard entities..."

# Function to create entity with error handling
create_entity() {
    local entity_id="$1"
    local name="$2"
    local icon="$3"
    local initial_state="$4"
    local attributes="$5"
    
    local payload="{
        \"state\": \"${initial_state}\",
        \"attributes\": {
            \"friendly_name\": \"${name}\",
            \"icon\": \"${icon}\",
            ${attributes}
        }
    }"
    
    local response
    response=$(curl -s -w "%{http_code}" -X POST \
        -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
        -H "Content-Type: application/json" \
        -d "${payload}" \
        "http://supervisor/core/api/states/${entity_id}" 2>/dev/null)
    
    local http_code="${response: -3}"
    if [[ "$http_code" =~ ^2[0-9][0-9]$ ]]; then
        bashio::log.debug "✅ Created entity: ${entity_id}"
        return 0
    else
        bashio::log.warning "⚠️ Failed to create entity: ${entity_id} (HTTP: ${http_code})"
        return 1
    fi
}

# Create dashboard entities
bashio::log.info "🎯 Creating status tiles..."

create_entity "sensor.actual_budget_system_status" \
    "System Status" \
    "mdi:calculator-variant" \
    "running" \
    '"description": "Overall system health", "tooltip": "Green=Running, Yellow=Issues, Red=Failed", "version": "2.0.1"'

create_entity "sensor.actual_budget_api_status" \
    "API Server" \
    "mdi:api" \
    "checking" \
    '"description": "Node.js API server status", "tooltip": "Connect to your existing API server", "external_url": "'${EXTERNAL_API_URL}'"'

create_entity "sensor.actual_budget_server_status" \
    "Actual Budget Server" \
    "mdi:server" \
    "checking" \
    '"description": "Official Actual Budget server", "tooltip": "Connect to your existing server", "external_url": "'${EXTERNAL_ACTUAL_URL}'"'

create_entity "sensor.actual_budget_last_sync" \
    "Last Sync" \
    "mdi:sync" \
    "never" \
    '"description": "Last synchronization time", "tooltip": "When transactions were last synced to Xero"'

bashio::log.info "📈 Creating statistics tiles..."

create_entity "sensor.actual_budget_transactions_processed" \
    "Transactions Processed" \
    "mdi:counter" \
    "0" \
    '"unit_of_measurement": "transactions", "description": "Total processed", "tooltip": "🔸 PLACEHOLDER: Shows test data until connected to API server"'

create_entity "sensor.actual_budget_successful_imports" \
    "Successful Imports" \
    "mdi:check-circle" \
    "0" \
    '"unit_of_measurement": "transactions", "description": "Successfully imported to Xero", "tooltip": "🔸 PLACEHOLDER: Shows test data until connected to API server"'

create_entity "sensor.actual_budget_failed_transactions" \
    "Failed Transactions" \
    "mdi:alert-circle" \
    "0" \
    '"unit_of_measurement": "transactions", "description": "Failed imports", "tooltip": "🔸 PLACEHOLDER: Shows test data until connected to API server"'

create_entity "sensor.actual_budget_pending_mappings" \
    "Pending Mappings" \
    "mdi:map-marker-question" \
    "0" \
    '"unit_of_measurement": "mappings", "description": "Awaiting category mapping", "tooltip": "🔸 PLACEHOLDER: Shows test data until connected to API server"'

bashio::log.info "🎉 Native HA dashboard integration complete!"
bashio::log.info "📱 Look for 'Actual Budget' in your Home Assistant sidebar"
bashio::log.info "📊 Dashboard tiles created with placeholder data"

# Show configuration info
if [[ -n "${EXTERNAL_API_URL}" ]]; then
    bashio::log.info "🔗 API Server URL: ${EXTERNAL_API_URL}"
else
    bashio::log.info "⚠️ No API server URL configured - configure external_api_url"
fi

if [[ -n "${EXTERNAL_ACTUAL_URL}" ]]; then
    bashio::log.info "🔗 Actual Budget URL: ${EXTERNAL_ACTUAL_URL}"
else
    bashio::log.info "⚠️ No Actual Budget URL configured - configure external_actual_url"
fi

# Function to update entity state with error handling
update_entity_state() {
    local entity_id="$1"
    local state="$2"
    local attributes="$3"
    
    local payload="{\"state\": \"${state}\", \"attributes\": ${attributes}}"
    
    curl -s -X POST \
        -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
        -H "Content-Type: application/json" \
        -d "${payload}" \
        "http://supervisor/core/api/states/${entity_id}" > /dev/null 2>&1
}

# Function to check service health with timeout
check_service_health() {
    local service_name="$1"
    local url="$2"
    
    if [[ -z "${url}" ]]; then
        return 1
    fi
    
    # Try to connect to the service with timeout
    if timeout 10 curl -s -f "${url}/health" > /dev/null 2>&1 || timeout 10 curl -s -f "${url}" > /dev/null 2>&1; then
        return 0
    else
        return 1
    fi
}

# Cleanup function for graceful shutdown
cleanup() {
    bashio::log.info "🛑 Received shutdown signal, cleaning up..."
    exit 0
}

# Set up signal handlers
trap cleanup SIGTERM SIGINT

# Main service loop - keep the add-on running and update entities periodically
bashio::log.info "🔄 Starting monitoring loop..."
bashio::log.info "💡 Add-on will run continuously and update dashboard every 2 minutes"

# Initialize loop counter
loop_count=0

# Infinite loop to keep the service running
while true; do
    ((loop_count++))
    bashio::log.debug "🔄 Health check cycle ${loop_count}"
    
    # Check API server health
    if [[ -n "${EXTERNAL_API_URL}" ]]; then
        if check_service_health "API Server" "${EXTERNAL_API_URL}"; then
            update_entity_state "sensor.actual_budget_api_status" "online" \
                '{"friendly_name": "API Server", "icon": "mdi:api", "description": "Connected and responding", "last_check": "'$(date -Iseconds)'"}'
            bashio::log.debug "✅ API Server: online"
        else
            update_entity_state "sensor.actual_budget_api_status" "offline" \
                '{"friendly_name": "API Server", "icon": "mdi:api-off", "description": "Cannot connect to API server", "last_check": "'$(date -Iseconds)'"}'
            bashio::log.debug "❌ API Server: offline"
        fi
    else
        update_entity_state "sensor.actual_budget_api_status" "not_configured" \
            '{"friendly_name": "API Server", "icon": "mdi:api-off", "description": "No external API URL configured", "last_check": "'$(date -Iseconds)'"}'
    fi
    
    # Check Actual Budget server health
    if [[ -n "${EXTERNAL_ACTUAL_URL}" ]]; then
        if check_service_health "Actual Budget Server" "${EXTERNAL_ACTUAL_URL}"; then
            update_entity_state "sensor.actual_budget_server_status" "online" \
                '{"friendly_name": "Actual Budget Server", "icon": "mdi:server", "description": "Connected and responding", "last_check": "'$(date -Iseconds)'"}'
            bashio::log.debug "✅ Actual Budget Server: online"
        else
            update_entity_state "sensor.actual_budget_server_status" "offline" \
                '{"friendly_name": "Actual Budget Server", "icon": "mdi:server-off", "description": "Cannot connect to Actual Budget server", "last_check": "'$(date -Iseconds)'"}'
            bashio::log.debug "❌ Actual Budget Server: offline"
        fi
    else
        update_entity_state "sensor.actual_budget_server_status" "not_configured" \
            '{"friendly_name": "Actual Budget Server", "icon": "mdi:server-off", "description": "No external Actual Budget URL configured", "last_check": "'$(date -Iseconds)'"}'
    fi
    
    # Update system status with current timestamp
    update_entity_state "sensor.actual_budget_system_status" "running" \
        '{"friendly_name": "System Status", "icon": "mdi:calculator-variant", "description": "Native HA integration active", "last_update": "'$(date -Iseconds)'", "version": "2.0.2", "integration_type": "native_dashboard", "loop_count": '${loop_count}'}'
    
    # Log status every 10 cycles (20 minutes)
    if (( loop_count % 10 == 0 )); then
        bashio::log.info "💓 Add-on running healthy - completed ${loop_count} health check cycles"
    fi
    
    # Sleep for 2 minutes before next health check
    # Use sleep with explicit duration to ensure proper blocking
    sleep 120 || true
done

# This line should never be reached
bashio::log.error "❌ Monitoring loop exited unexpectedly!"
exit 1